# systemd/.config/systemd/user/checkmail.service
[Unit]
Description=check mail
RefuseManualStart=no
RefuseManualStop=yes

[Service]
Type=oneshot
ExecStart=/home/kef/.bin/checkmail.sh
